# supervisord.ini
[program:gunicorn]
command = python manage.py runserver 0.0.0.0:8010
stdout_logfile = /dev/stdout
stderr_logfile = /dev/stderr


[program:grpc]
command = python manage.py grpcrunserver --dev
stdout_logfile = /dev/stdout
stderr_logfile = /dev/stderr
