# supervisor/supervisor.conf
[program:gunicorn]
command = python manage.py runserver 0.0.0.0:8010


[program:grpc]
command = python manage.py grpcrunserver --dev
